# ors_graph_update_scripts/fetch_graphs.cron
#!/bin/sh

/opt/ors/data/fetch_graphs.sh $GRAPH_TYPE $SIBLING_WORKER $DOCKER_CONTAINER_NAME > "/opt/ors/data/fetch_graphs.out.$(date +%Y%m%d_%H%M%S)"
